# v5.0/scripts/oddsapi-snapshot-closing.service
[Unit]
Description=SportStock V5.0 — OddsAPI Closing Snapshot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=derek
Group=derek
WorkingDirectory=/home/derek/Documents/stock/v5.0

EnvironmentFile=-/home/derek/Documents/stock/v5.0/.env

ExecStart=/usr/bin/python3 scripts/odds_scheduler.py --provider oddsapi --mode closing --once --normalize
StandardOutput=append:/home/derek/Documents/stock/v5.0/logs/oddsapi_snapshot_closing.log
StandardError=append:/home/derek/Documents/stock/v5.0/logs/oddsapi_snapshot_closing.log
SyslogIdentifier=v5-oddsapi-snapshot-closing
TimeoutStartSec=600

[Install]
WantedBy=default.target
