# v5.0/scripts/odds-opening.timer
# ──────────────────────────────────────────────────────────
# SportStock V5.0 — Opening Odds Timer
#
# Fires daily at 05:05 UTC (12:05 AM EST) to capture
# opening odds lines for all sports.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp scripts/odds-opening.service ~/.config/systemd/user/
#   cp scripts/odds-opening.timer   ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable odds-opening.timer
#   systemctl --user start odds-opening.timer
#
# Check status:
#   systemctl --user list-timers
#   systemctl --user status odds-opening.timer
#   journalctl --user -u odds-opening.service
# ──────────────────────────────────────────────────────────

[Unit]
Description=SportStock V5.0 — Opening Odds Timer (daily 12:05 AM EST)

[Timer]
OnCalendar=*-*-* 05:05:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
