# v5.0/scripts/sgo-snapshot-closing.timer
[Unit]
Description=SportStock V5.0 — SGO Closing Odds Snapshot Timer (23:00 UTC daily)
Requires=sgo-snapshot-closing.service

[Timer]
# 23:00 UTC = 6:00 PM EST — before US evening tip-offs
OnCalendar=*-*-* 23:00:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
